# systemd/locksmithd.service
[Unit]
Description=Cluster reboot manager
After=etcd.service
After=update-engine.service
After=user-config.target
After=system-config.target

[Service]
ExecStart=/usr/lib/locksmith/locksmithd
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
